# roles/steam/templates/steam-move-screenshots.timer.j2
[Unit]
Description=Hourly move of steam screenshots to Pictures dir

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
